# docker/postfix/supervisord.conf
[supervisord]
nodaemon=true

[program:postfix]
command=/opt/postfix.sh